# support-chat/huggingface/supervisord.conf
[supervisord]
nodaemon=true
user=appuser
logfile=/var/log/supervisor/supervisord.log
pidfile=/tmp/supervisord.pid

[program:registry]
command=java -Xms128m -Xmx512m -Dquarkus.http.host=0.0.0.0 -Dquarkus.http.port=8080 -Djava.util.logging.manager=org.jboss.logmanager.LogManager -jar /opt/registry/quarkus-run.jar
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:support-chat]
command=java -Xms128m -Xmx512m -Dquarkus.http.host=0.0.0.0 -Dquarkus.http.port=7860 -Djava.util.logging.manager=org.jboss.logmanager.LogManager -jar /opt/support-chat/quarkus-run.jar
environment=REGISTRY_URL="http://localhost:8080",CORS_ORIGINS="*",HTTP_PORT="7860"
autostart=true
autorestart=true
startsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
